# systemd/user/nanoclaw-audit-audio.service
[Unit]
Description=NanoClaw §201 filesystem audio audit (REQ-QUAL-04 / REQ-LEGAL-03)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-%h/.config/nanoclaw/audit.env
ExecStart=%h/nanoclaw/scripts/audit-audio.sh
StandardOutput=journal
StandardError=journal
